# systemd/system/system-walinuxagent.extensions-Microsoft.OSTCExtensions.VMAccessForLinux_1.5.4.slice

[Unit]
Description=Slice for extension Microsoft.OSTCExtensions.VMAccessForLinux-1.5.4
DefaultDependencies=no
Before=slices.target
Requires=system-walinuxagent.extensions.slice
After=system-walinuxagent.extensions.slice